# deploy/tracker.socket
[Unit]
Description=gunicorn socket para o Trankeador

[Socket]
ListenStream=/run/tracker.sock

[Install]
WantedBy=sockets.target
